# files/upstart/init/nova-compute.conf
description "Nova Compute server"
author "Yun Mao <yunmao@gmail.com>"

start on (filesystem and net-device-up IFACE!=lo)
stop on runlevel [016]

chdir %DIR%

pre-start script
        mkdir -p /var/log/nova
        chown %USER%:root /var/log/nova/
        modprobe nbd
end script

#respawn

exec su -c "%DIR%/bin/nova-compute --logdir=%LOGDIR%/nova" %USER%
